# packer/base/root/etc/cron.hourly/logrotate
#!/bin/sh
exec logrotate /etc/logrotate.conf
